# contrib/init/ubicoind.conf
description "Ubicoin Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env UBICOIND_BIN="/usr/bin/ubicoind"
env UBICOIND_USER="ubicoin"
env UBICOIND_GROUP="ubicoin"
env UBICOIND_PIDDIR="/var/run/ubicoind"
# upstart can't handle variables constructed with other variables
env UBICOIND_PIDFILE="/var/run/ubicoind/ubicoind.pid"
env UBICOIND_CONFIGFILE="/etc/ubicoin/ubicoin.conf"
env UBICOIND_DATADIR="/var/lib/ubicoind"

expect fork

respawn
respawn limit 5 120
kill timeout 60

pre-start script
    # this will catch non-existent config files
    # ubicoind will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$UBICOIND_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run ubicoind."
        echo "The setting must appear in $UBICOIND_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following "
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"Ubicoin Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$UBICOIND_PIDDIR"
    chmod 0755 "$UBICOIND_PIDDIR"
    chown $UBICOIND_USER:$UBICOIND_GROUP "$UBICOIND_PIDDIR"
    chown $UBICOIND_USER:$UBICOIND_GROUP "$UBICOIND_CONFIGFILE"
    chmod 0660 "$UBICOIND_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$UBICOIND_PIDFILE" \
    --chuid $UBICOIND_USER:$UBICOIND_GROUP \
    --exec "$UBICOIND_BIN" \
    -- \
    -pid="$UBICOIND_PIDFILE" \
    -conf="$UBICOIND_CONFIGFILE" \
    -datadir="$UBICOIND_DATADIR" \
    -disablewallet \
    -daemon

